# .config/autostart/screenlock.desktop
[Desktop Entry]
Name=Screen Lock
Type=Application
Exec=xss-lock -- xscreensaver-command -lock
