# .config/autostart/urxvtd.desktop
[Desktop Entry]
Type=Application
Name=Rxvt Unicode Daemon
Icon=terminal
Exec=urxvtd -q -o -f
StartupNotify=true
Categories=Terminal;
